# systemd/readsbstats-updater.service
[Unit]
Description=readsbstats aircraft/airline database updater
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=readsbstats
Group=readsbstats
# v2.12.2: stop the collector before the updater body runs and restart it
# after, so the weekly timer path can't collide with concurrent writers on
# the v2.11.7+ single-transaction aircraft_db swap. `+` runs as root,
# bypassing User=/hardening for these two systemctl calls only — User=
# readsbstats cannot manage system units without polkit. ExecStopPost
# always runs, so a crashed/timed-out updater still restores the collector.
ExecStartPre=+/bin/systemctl stop readsbstats-collector.service
ExecStart=/opt/readsbstats/venv/bin/python -m readsbstats.db_updater
ExecStopPost=+/bin/systemctl start readsbstats-collector.service
WorkingDirectory=/opt/readsbstats
# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/mnt/ext/readsbstats
# v2.12.2: 300 → 900 — the v2.11.7+ single-transaction swap holds
# BEGIN IMMEDIATE for the full CREATE + chunked INSERT + rename sequence;
# on Pi-4 SD/USB that's 3–5 min for 620 k rows. 900 s gives ~3× headroom.
TimeoutStartSec=900

# Audit-13 A13-046: kernel + namespace hardening.
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=

StandardOutput=journal
StandardError=journal
SyslogIdentifier=readsbstats-updater
