# systemd/readsbstats-collector.service
[Unit]
Description=readsbstats ADS-B History Collector
After=network.target readsb.service
Wants=readsb.service
OnFailure=notify-telegram@%n.service
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=notify
WatchdogSec=60
User=readsbstats
Group=readsbstats
WorkingDirectory=/opt/readsbstats
ExecStart=/opt/readsbstats/venv/bin/python -m readsbstats.collector
Restart=on-failure
RestartSec=10
TimeoutStartSec=300

# Be polite to the Pi — readsb/tar1090 are the primary tenants
CPUQuota=15%
MemoryMax=192M
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/mnt/ext/readsbstats
ReadOnlyPaths=/run/readsb

# Audit-13 A13-046: kernel + namespace hardening. Targets systemd-analyze
# security score < 5. MemoryDenyWriteExecute is intentionally omitted —
# DuckDB has known issues with MDWE on some builds. Verify on the Pi
# with `systemd-analyze security readsbstats-collector.service` and add
# `MemoryDenyWriteExecute=yes` once verified safe.
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=

# Logging goes to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=readsbstats-collector

[Install]
WantedBy=multi-user.target
