# systemd/readsbstats-dbcheck-full.service
[Unit]
Description=readsbstats monthly DB integrity_check
OnFailure=notify-telegram@%n.service
# Never run the heavy integrity_check concurrently with the weekly quick_check
# (they collide on the 1st Sunday). Conflicts is symmetric: starting this stops
# a still-running quick check — full supersedes it, so no coverage is lost.
Conflicts=readsbstats-dbcheck.service

[Service]
Type=oneshot
User=readsbstats
Group=readsbstats
ExecStart=/opt/readsbstats/venv/bin/python /opt/readsbstats/scripts/check_db.py --mode full
WorkingDirectory=/opt/readsbstats
# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/mnt/ext/readsbstats
# integrity_check is far heavier than quick_check (it cross-checks indexes vs
# tables) and reads the whole DB off the USB disk. Scaled up alongside the
# quick-check bump: on a multi-GB DB with a busy collector this can run many
# minutes, so 600 s could trip under load. 1800 s (monthly, off-peak) is safe.
TimeoutStartSec=1800

# Audit-13 A13-046: kernel + namespace hardening.
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=

StandardOutput=journal
StandardError=journal
SyslogIdentifier=readsbstats-dbcheck-full
