# systemd/readsbstats-dbcheck.service
[Unit]
Description=readsbstats weekly DB quick_check
OnFailure=notify-telegram@%n.service

[Service]
Type=oneshot
User=readsbstats
Group=readsbstats
ExecStart=/opt/readsbstats/venv/bin/python /opt/readsbstats/scripts/check_db.py --mode quick
WorkingDirectory=/opt/readsbstats
# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/mnt/ext/readsbstats
# quick_check is read-only but reads every page off the (USB) disk while the
# collector writes; a read-only WAL reader also walks the -wal frames, so wall
# time scales with DB + WAL size and disk contention, not CPU. Measured ~60 s
# standalone on a 1.2 GB DB with a 119 MB WAL, and >120 s during a deploy/reboot
# I/O storm. 120 s was too tight (false 'timeout' failures + Telegram alerts);
# 600 s leaves headroom as the DB grows.
TimeoutStartSec=600

# Audit-13 A13-046: kernel + namespace hardening.
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
CapabilityBoundingSet=

StandardOutput=journal
StandardError=journal
SyslogIdentifier=readsbstats-dbcheck
