# systemd/readsbstats-dbcheck-full.service
[Unit]
Description=readsbstats monthly DB integrity_check
OnFailure=notify-telegram@%n.service

[Service]
Type=oneshot
User=readsbstats
Group=readsbstats
ExecStart=/opt/readsbstats/venv/bin/python /opt/readsbstats/scripts/check_db.py --mode full
WorkingDirectory=/opt/readsbstats
# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/mnt/ext/readsbstats
TimeoutStartSec=600

StandardOutput=journal
StandardError=journal
SyslogIdentifier=readsbstats-dbcheck-full
